# big-preload/etc/xdg/autostart/big-preload-user.desktop
[Desktop Entry]
Name=Preload user desktop
Exec=big-preload-user-gen-list
Terminal=false
Type=Application
StartupNotify=false
Categories=System
