# big-preload/DEBIAN/postinst
#!/bin/sh

systemctl enable big-preload

exit 0
