# big-preload/DEBIAN/prerm
#!/bin/sh

systemctl disable big-preload

exit 0
